# dev-util/tuxdroid/files/tuxdroid.rc6
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $
 
depend() {
        need alsasound
}
 
start() {
        ebegin "Starting Tux HTTP Server"
        start-stop-daemon --start -p /var/run/tuxhttpserver.pid --exec /usr/share/tuxdroid/tuxhttpserver/tuxhttpserver.py --name python
        eend $?
}
 
stop() {
        ebegin "Stopping Tux HTTP Server"
        start-stop-daemon --stop --quiet -p /var/run/tuxhttpserver.pid 
        eend $?
}
